# scripts/kaia-lockdown.service
# IMPORTANT: Set KAIA_PROJECT_DIR to the absolute path of your Kaia repo.
# Edit the Environment= line below before running: systemctl enable --now <service>
[Unit]
Description=Kaia Emergency Network Lockdown
DefaultDependencies=no
Before=network.target shutdown.target

[Service]
Type=oneshot
Environment=KAIA_PROJECT_DIR=/home/ekco/github/Kaia
ExecStart=/bin/bash ${KAIA_PROJECT_DIR}/scripts/kaia-lockdown.sh
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
