# scripts/kaia-policy-gate.service
# IMPORTANT: Set KAIA_PROJECT_DIR to the absolute path of your Kaia repo.
# Edit the Environment= line below before running: systemctl enable --now <service>
[Unit]
Description=Kaia Out-of-Process Policy Gate Daemon
After=network.target postgresql.service

[Service]
Type=simple
Environment=KAIA_PROJECT_DIR=/home/ekco/github/Kaia
ExecStart=/usr/bin/python3 ${KAIA_PROJECT_DIR}/security/policy_gate.py
WorkingDirectory=${KAIA_PROJECT_DIR}
Restart=always
RestartSec=1
RuntimeDirectory=kaiacord
User=root
Group=kaiacord
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ReadWritePaths=/run/kaiacord
CapabilityBoundingSet=CAP_NET_RAW CAP_SYS_ADMIN
AmbientCapabilities=CAP_NET_RAW CAP_SYS_ADMIN

[Install]
WantedBy=multi-user.target
